# systemd/ggnet-worker.service
[Unit]
Description=GGnet Image Conversion Worker
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=ggnet
Group=ggnet
WorkingDirectory=/opt/ggnet/backend
ExecStart=/opt/ggnet/venv/bin/python -m worker.convert
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ggnet-worker

# Environment variables
Environment=PYTHONPATH=/opt/ggnet/backend
Environment=GGNET_ENV=production

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ggnet/storage /opt/ggnet/logs

# Resource limits
LimitNOFILE=65536
MemoryMax=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target